# systemd/k8s-agent-tools-server-backup.timer
[Unit]
Description=Daily backup for K8s Agent Tools Server

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
